# deploy/systemd/ethusdc-pivot-bot-dry-run.user.service
[Unit]
Description=ETHUSDC Pivot Bot WS-only dry-run
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/dev/ethusdc-pivot-bot
EnvironmentFile=/home/dev/ethusdc-pivot-bot/.env
Environment=DRY_RUN=true
Environment=LIVE_TRADING=false
Environment=PUBLIC_MARKET_DRY_RUN=true
Environment=PUBLIC_MARKET_WS_ONLY=true
Environment=ALLOW_CACHED_EXCHANGE_FILTERS_IN_DRY_RUN=true
Environment=BINANCE_SYMBOL=ETHUSDC
Environment=BINANCE_INTERVAL=15m
Environment=EXIT_AFTER_BOUNDED_RUNTIME=true
Environment=PHASE3B_BOUNDED_RUNTIME_SECONDS=3600
Environment=BINANCE_API_KEY=
Environment=BINANCE_API_SECRET=
ExecStart=/home/dev/ethusdc-pivot-bot/.venv/bin/python -m bot.main
Restart=on-failure
RestartSec=10
RuntimeMaxSec=3700
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=default.target
