# deploy/systemd/ethusdc-pivot-bot-live-monitor.user.timer
[Unit]
Description=Run ETHUSDC Pivot Bot read-only live monitor every 15 minutes

[Timer]
OnBootSec=2min
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
