# deploy/systemd/xrpusdc-pivot-bot-live-monitor.user.service
[Unit]
Description=XRPUSDC Pivot Bot read-only live monitor status writer

[Service]
Type=oneshot
WorkingDirectory=%h/ethusdc-pivot-bot
ExecStart=%h/ethusdc-pivot-bot/scripts/run_live_monitor_status_xrpusdc_1h.sh
